# etc/systemd/system/open-fprintd-resume.service.d/local.conf
#[Unit]
#PropagatesReloadTo=python3-validity.service

[Service]
Restart=on-failure
ExecStartPre=systemctl restart python3-validity.service
